# src/supervisor/conf.d/celery_beat.conf
[program:celery_beat]
command=celery -A moodtracker beat
user=root
numprocs=1
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile_maxbytes=1MB
stderr_logfile_maxbytes=1MB
stdout_logfile=/var/log/celery_beat.log
stderr_logfile=/var/log/celery_beat_err.log